# battery.timer
[Unit]
Description="Timer wich runs battery every minute"

[Timer]
OnUnitActiveSec=60
Unit=battery.service

[Install]
WantedBy=default.target

